# supervisord.conf
; RUBY tracing client etl - pylonsd

[program:pylonsd]
command = /home/big-dipper/start.sh

[program:ruby-etl]
command =  bash -c 'tail -f /tmp/trace/trace.fifo | ruby /home/big-dipper/trace.rb'

[eventlistener:processes]
command=bash -c "printf 'READY\n' && while read line; do kill -SIGQUIT $PPID; done < /dev/stdin"
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL

[supervisord]
#logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
#logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
#logfile_backups=10                              ; number of backed up logfiles
loglevel=trace                                 ; info, debug, warn, trace
pidfile=/var/log/supervisord/supervisord.pid                ; pidfile location
nodaemon=yes                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=big-dipper                                      ; default user
#childlogdir=/var/log/supervisord/               ; where child log files will live
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

